# .config/autostart/pigdin.desktop
[Desktop Entry]
Type=Application
Exec=pidgin
