# .config/autostart/volume.desktop
[Desktop Entry]
Type=Application
Exec=volumeicon
